# config/supervisor-cron.ini
[program:cron]
command=crond -f -L %(ENV_CRON_STDOUT_FILE)s
priority=2
autostart=true
autorestart=true
stdout_logfile=%(ENV_CRON_STDOUT_FILE)s
stderr_logfile=%(ENV_CRON_STDERR_FILE)s

